# scripts/max-reflect.service.template
[Unit]
Description=Max Reflection Task — dispatches cog-reflect via the running daemon (requires max.service)
After=max.service
Requires=max.service

[Service]
Type=oneshot
ExecStart=__MAX_BIN__ reflect
WorkingDirectory=__HOME__
Environment=PATH=__PATH__
Environment=HOME=__HOME__
EnvironmentFile=-__ENV_FILE__
StandardOutput=journal
StandardError=journal
SyslogIdentifier=max-reflect
TimeoutStartSec=300
